# alerts/rules/imported_from_Persistence_Via_Cron_Files_6c4e2f43-d94d-4ead-b64d-97e53fa2bd05.json
{
  "rule_name": "imported_from_Persistence_Via_Cron_Files_6c4e2f43-d94d-4ead-b64d-97e53fa2bd05",
  "rule_author": "siembol",
  "rule_version": 1,
  "rule_description": "Detects creation of cron file or files in Cron directories which could indicates potential persistence.",
  "source_type": "*",
  "matchers": [
    {
      "is_enabled": true,
      "matcher_type": "COMPOSITE_OR",
      "is_negated": false,
      "matchers": [
        {
          "is_enabled": true,
          "matcher_type": "REGEX_MATCH",
          "is_negated": false,
          "field": "TargetFilename",
          "data": "^\\Q/etc/cron.d/\\E.*|^\\Q/etc/cron.daily/\\E.*|^\\Q/etc/cron.hourly/\\E.*|^\\Q/etc/cron.monthly/\\E.*|^\\Q/etc/cron.weekly/\\E.*|^\\Q/var/spool/cron/crontabs/\\E.*"
        },
        {
          "is_enabled": true,
          "matcher_type": "REGEX_MATCH",
          "is_negated": false,
          "field": "TargetFilename",
          "data": ".*\\Q/etc/cron.allow\\E.*|.*\\Q/etc/cron.deny\\E.*|.*\\Q/etc/crontab\\E.*"
        }
      ]
    }
  ],
  "tags": [
    {
      "tag_name": "sigma_id",
      "tag_value": "6c4e2f43-d94d-4ead-b64d-97e53fa2bd05"
    },
    {
      "tag_name": "sigma_tags",
      "tag_value": "[attack.persistence, attack.t1053.003]"
    },
    {
      "tag_name": "sigma_fp",
      "tag_value": "[Any legitimate cron file.]"
    }
  ],
  "rule_protection": {
    "max_per_hour": 30,
    "max_per_day": 100
  }
}